# systemd/system/fake_buttons@.service
[Unit]
Description=Fake Input Server

[Service]
ExecStart=-/usr/sbin/fake_buttons
StandardInput=socket
User=www-data
Group=www-data
