# systemd/system/results.service
[Unit]
Description=Static Results Generator
After=timing.service
ConditionPathExists=/etc/timing/timing.conf
ConditionPathExists=/data

[Service]
Type=simple
IOSchedulingClass=best-effort
IOSchedulingPriority=3
WorkingDirectory=/data
User=www-data
Group=www-data
ExecStart=/usr/bin/results_to_static
Restart=always

[Install]
WantedBy=multi-user.target
