# etc/systemd/system/hostapd@.service.d/override.conf
[Unit]
Requisite=sys-subsystem-net-devices-%i.device