# systemd/system/pigpiod.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/pigpiod -l -n 127.0.0.1
# Change sample rate to 1 micro second from default of 5
#ExecStart=/usr/bin/pigpiod -l -n 127.0.0.1 -s 1